# build/oracle.conf
[program:start-script]
autorestart = false
startsecs = 0
user = root
command = /root/setup.sh
umask = 000