# build/system/oracle.service
[Unit]
Description=Helium Oracle

[Service]
Type=simple
ExecStart=/root/submit.sh